# quadlet/game-dev.container
# game-dev.container — example workload (Vite + Phaser). Replace with your own
# project; the server itself is "opencode-dev-server".
[Unit]
Description=Vite + Phaser game dev server
After=network-online.target dev-network.service
Wants=network-online.target

[Container]
Image=docker.io/node:22-alpine
# Pin by digest in production: Image=docker.io/node@sha256:...
ContainerName=game-dev
Network=dev.network

# Bind only to WireGuard interface
PublishPort=10.44.0.1:5173:5173

# Mount repository from persistent state
Volume=/mnt/state/repos:/repos:Z
# Mount npm cache
Volume=/mnt/state/cache/npm:/root/.npm:Z

WorkingDir=/repos/game
Exec=sh -c "npm install && npx vite --host 0.0.0.0 --port 5173"

Environment=NODE_ENV=development

[Service]
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target default.target
