# config/butane/files/systemd/base/wg-quick@wg0.service.d/after-wg-setup.conf
[Unit]
After=wg-setup.service
Requires=wg-setup.service
