# config/butane/files/etc/containers/systemd/users/1000/opencode.container
[Unit]
Description=OpenCode headless server (VPN-only)
After=network-online.target
Wants=network-online.target

[Container]
Image=opencode.build
ContainerName=opencode
PublishPort=10.44.0.1:4096:4096
Volume=/mnt/state/opencode:/data:Z
Environment=XDG_DATA_HOME=/data
Environment=XDG_CONFIG_HOME=/data/config
Environment=HOME=/data
Environment=POCKETBASTION=1
Volume=/mnt/state/repos:/repos:Z
WorkingDir=/repos
EnvironmentFile=/mnt/state/secrets/opencode.env
Exec=serve --hostname 0.0.0.0 --port 4096

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
